# services/pageboard.service
# This service is actually a systemd target, but we are using a service since
# targets cannot be reloaded.

[Unit]
Description=Control %p
Requires=%p@fr.service
Requires=%p@io.service

[Service]
Type=oneshot
ExecStart=/bin/true
ExecReload=/bin/true
RemainAfterExit=on

[Install]
WantedBy=default.target
